# cookbooks/pookio/templates/default/pdf_saver_server_initd.sh.erb
#!/bin/bash
# Initialize script for thin daemons
### BEGIN INIT INFO
# Provides:          pdf_saver_server
# Required-Start:    $local_fs $remote_fs $network
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: pdf_saver_server initscript
# Description:       pdf_saver_server part of the pookio server family
### END INIT INFO
# prototypical init script
# Do NOT "set -e"

USER=pookio
PB4US_ROOT="<%= node['pookio']['root'] %>"
PDF_SAVER_SERVER="${PB4US_ROOT}/pb_server/current/script/pdf_saver_server"
SECRETS="${PB4US_ROOT}/pb_server/shared/secrets.sh"
export RBENV_ROOT="/usr/local/rbenv"
RBENV="/etc/profile.d/rbenv.sh"
export RACK_ENV="<%= node['pookio']['rack_env'] %>"
CMD="bundle exec $PDF_SAVER_SERVER"

. "$RBENV"
cd "$PB4US_ROOT/pb_server/current"


case "$1" in
  start)
    . $SECRETS && su $USER -c "PATH=$PATH $CMD start"
    ;;
  stop)
    . $SECRETS && su $USER -c "PATH=$PATH $CMD stop"
    ;;
  restart)
    # . $SECRETS && su $USER -c "PATH=$PATH $CMD restart"
    . $SECRETS && su $USER -c "PATH=$PATH $CMD stop"
    . $SECRETS && su $USER -c "PATH=$PATH $CMD start"
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}" >&2
    exit 3
    ;;
esac
